# crontab
#================ LIGAR QUARTA =================
#19:30
25 19 * * 3 http://192.168.2.220:1000/?acao=00ON
00 19 * * 3 http://192.168.2.220:1000/?acao=01ON
02 19 * * 3 http://192.168.2.220:1000/?acao=02ON
04 19 * * 3 http://192.168.2.220:1000/?acao=03ON
10 19 * * 3 http://192.168.2.220:1000/?acao=04ON
15 19 * * 3 http://192.168.2.220:1000/?acao=05ON
18 19 * * 3 http://192.168.2.220:1000/?acao=06ON
#DESLILIGAR QUARTA 21:30
10 21 * * 3 http://192.168.2.220:1000/?acao=00OF
12 21 * * 3 http://192.168.2.220:1000/?acao=01OF
14 21 * * 3 http://192.168.2.220:1000/?acao=02OF
18 21 * * 3 http://192.168.2.220:1000/?acao=03OF
20 21 * * 3 http://192.168.2.220:1000/?acao=04OF
22 21 * * 3 http://192.168.2.220:1000/?acao=05OF
24 21 * * 3 http://192.168.2.220:1000/?acao=06OF

#================ LIGAR DOMINGO =================
#LIGAR DOMINGO 9:30 
25 9 * * 0 http://192.168.2.220:1000/?acao=00ON
00 9 * * 0 http://192.168.2.220:1000/?acao=01ON
02 9 * * 0 http://192.168.2.220:1000/?acao=02ON
04 9 * * 0 http://192.168.2.220:1000/?acao=03ON
10 9 * * 0 http://192.168.2.220:1000/?acao=04ON
15 9 * * 0 http://192.168.2.220:1000/?acao=05ON
18 9 * * 0 http://192.168.2.220:1000/?acao=06ON

#DESLILIGAR DOMINGO 12:00
35 11 * * 0 http://192.168.2.220:1000/?acao=00OF
40 11 * * 0 http://192.168.2.220:1000/?acao=01OF
45 11 * * 0 http://192.168.2.220:1000/?acao=02OF
50 11 * * 0 http://192.168.2.220:1000/?acao=03OF
42 11 * * 0 http://192.168.2.220:1000/?acao=04OF
49 11 * * 0 http://192.168.2.220:1000/?acao=05OF
38 11 * * 0 http://192.168.2.220:1000/?acao=06OF

#LIGAR DOMINGO 16:00 
10 16 * * 0 http://192.168.2.220:1000/?acao=00ON
15 16 * * 0 http://192.168.2.220:1000/?acao=01ON
20 16 * * 0 http://192.168.2.220:1000/?acao=02ON
25 16 * * 0 http://192.168.2.220:1000/?acao=03ON
30 16 * * 0 http://192.168.2.220:1000/?acao=04ON
35 16 * * 0 http://192.168.2.220:1000/?acao=05ON
40 16 * * 0 http://192.168.2.220:1000/?acao=06ON

#DESLILIGAR DOMINGO 21:00
35 21 * * 0 http://192.168.2.220:1000/?acao=00OF
40 21 * * 0 http://192.168.2.220:1000/?acao=01OF
45 21 * * 0 http://192.168.2.220:1000/?acao=02OF
50 21 * * 0 http://192.168.2.220:1000/?acao=03OF
42 21 * * 0 http://192.168.2.220:1000/?acao=04OF
49 21 * * 0 http://192.168.2.220:1000/?acao=05OF
38 21 * * 0 http://192.168.2.220:1000/?acao=06OF
